# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:jupyter]
command=jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --allow-root --NotebookApp.token='' --NotebookApp.password=''
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/jupyter_stdout.log
stderr_logfile=/var/log/supervisor/jupyter_stderr.log

[program:worker]
command=python3 /workspace/FFServer/worker.py
directory=/workspace/FFServer
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/worker_stdout.log
stderr_logfile=/var/log/supervisor/worker_stderr.log

[program:streamlit]
command=streamlit run app.py --server.address 0.0.0.0
directory=/workspace/FFServer
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/streamlit_stdout.log
stderr_logfile=/var/log/supervisor/streamlit_stderr.log
